# deploy/muneral-kb-sync.service
[Unit]
Description=Muneral to Scrutator incremental graph sync
Documentation=https://github.com/Arcanada-one/scrutator
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=15min
StartLimitBurst=3
ConditionPathExists=/etc/muneral-kb-sync/muneral-db-dsn
ConditionPathExists=/etc/muneral-kb-sync/ltm-writer-token

[Service]
Type=oneshot
User=muneral-kb-sync
Group=muneral-kb-sync
LoadCredential=muneral-db-dsn:/etc/muneral-kb-sync/muneral-db-dsn
LoadCredential=ltm-writer-token:/etc/muneral-kb-sync/ltm-writer-token
StateDirectory=muneral-kb-sync
StateDirectoryMode=0750
WorkingDirectory=/opt/muneral-kb-sync/current
ExecStart=/opt/muneral-kb-sync/current/bin/muneral-kb-sync --incremental --timer
TimeoutStartSec=10min
Restart=on-failure
RestartSec=30s
UMask=0077

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
SystemCallArchitectures=native
ReadWritePaths=/var/lib/muneral-kb-sync
